# deploy/linux/autolab-client.service
[Unit]
Description=AutoLab hardware push client
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=__WORKDIR__
EnvironmentFile=__WORKDIR__/.env
ExecStart=/usr/bin/env bash -lc 'docker run --rm --name autolab-client --env-file __WORKDIR__/.env autolab-client:latest'
ExecStop=/usr/bin/env bash -lc 'docker stop autolab-client || true'
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
